# frpcs_docker_crowsec/frpcs.conf
[supervisord]
nodaemon=true

[program:frps]
command=/frp/frps -c /frp/frps.ini 

[program:frpc]
command=/frp/frpc -c /frp/frpc.ini
